# deploy/poly-bot.service
[Unit]
Description=Polymarket Copy-Trading Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/poly_copy
ExecStart=/usr/bin/python3 -u bot.py --interval 20
Restart=always
RestartSec=30
CPUQuota=50%
MemoryMax=400M
StandardOutput=append:/opt/poly_copy/bot.log
StandardError=append:/opt/poly_copy/bot.log

[Install]
WantedBy=multi-user.target
