# deploy/poly-bot.service
[Unit]
Description=Polymarket Copy-Trading Bot (v1.3)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/poly_copy
ExecStart=/usr/bin/python3 -u bot.py
Restart=always
RestartSec=10
StandardOutput=append:/opt/poly_copy/bot.log
StandardError=append:/opt/poly_copy/bot.log

[Install]
WantedBy=multi-user.target
